# notes-snapshot.service
[Unit]
Description=Snapshot notes and push to origin

[Service]
ExecStart=/usr/local/bin/notes-snapshot.sh
