# notes-snapshot.timer
[Unit]
Description=Snapshot notes hourly

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
